# nskm.service
[Unit]
Description=Non-Sucking Key Mapper

[Service]
ExecStart=/usr/bin/nskm /dev/input/event%i

[Install]
WantedBy=getty.target multi-user.target
