# nskm.service
[Unit]
Description=Non-Sucking Key Mapper

[Service]
ExecStart=/usr/bin/nskm /dev/input/by-path/platform-i8042-serio-0-event-kbd
